# samba-backup.service
[Unit]
Description=Бэкап Samba-шар (samba-admin panel)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/samba-backup.sh
